# bmtl-device.service
[Unit]
Description=BMTL Device MQTT Client Daemon
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/bmtl-device/current/venv/bin/python

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bmtl-device/current
Environment=PYTHONPATH=/opt/bmtl-device/current
# Ensure system binaries like git and gphoto2 are available
Environment=PATH=/opt/bmtl-device/current/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/bmtl-device/launch-device.sh
Restart=always
RestartSec=15
StartLimitBurst=10
StartLimitIntervalSec=300
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=false
ProtectHome=yes
ReadWritePaths=/opt/bmtl-device /etc/bmtl-device

[Install]
WantedBy=multi-user.target
