# bmtl-camera.service
[Unit]
Description=BMTL Camera Control Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bmtl-device
Environment=PYTHONPATH=/opt/bmtl-device
ExecStart=/opt/bmtl-device/venv/bin/python /opt/bmtl-device/camera_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Minimal security settings to avoid namespace conflicts
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target