# bmtl-camera.service
[Unit]
Description=BMTL Camera Control Daemon
After=network-online.target bmtl-device.service
Wants=network-online.target
PartOf=bmtl-device.service
ConditionPathExists=/opt/bmtl-device/current/venv/bin/python

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bmtl-device/current
Environment=PYTHONPATH=/opt/bmtl-device/current
ExecStart=/opt/bmtl-device/launch-camera.sh
Restart=always
RestartSec=15
StartLimitBurst=10
StartLimitIntervalSec=300
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal

# Minimal security settings to avoid namespace conflicts
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
