# bmtl-camera.service
[Unit]
Description=BMTL Camera Control Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bmtl-device
Environment=PYTHONPATH=/opt/bmtl-device
ExecStartPre=/bin/mkdir -p /tmp/bmtl-config
ExecStart=/opt/bmtl-device/venv/bin/python /opt/bmtl-device/camera_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bmtl-device /tmp/bmtl-config
DeviceAllow=/dev/bus/usb

[Install]
WantedBy=multi-user.target