# bmtl-device.service
[Unit]
Description=BMTL Device MQTT Client Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bmtl-device/current
Environment=PATH=/opt/bmtl-device/current/venv/bin
ExecStart=/opt/bmtl-device/current/venv/bin/python /opt/bmtl-device/current/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=false
ProtectHome=yes
ReadWritePaths=/opt/bmtl-device /etc/bmtl-device

[Install]
WantedBy=multi-user.target